# mariadb/rootfs/etc/periodic/monthly/monthly_export
#!/command/with-contenv bashio
# ==============================================================================
# Export homeassistant database periodically
# ==============================================================================

if bashio::config.equals "purge_statistics.period" "monthly"; then
    purge-homeassistant-statistics-database
fi

if bashio::config.true "retention.enabled" && bashio::config.equals "retention.periodic" "monthly"; then
    export-databases
fi
