# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core-post/finish
#!/command/with-contenv bashio
# ==============================================================================
# Post-start finalization of  MariaDB service
# ==============================================================================

# Delete service information from the Supervisor
bashio::log.info "Deleting service information from Home Assistant."
bashio::services.delete "mysql"

# Export database content
export-homeassistant-database
