# mariadb/rootfs/etc/periodic/daily/daily_export
#!/command/with-contenv bashio
# ==============================================================================
# Export all databases periodically
# ==============================================================================

if bashio::config.equals "purge_statistics.period" "daily"; then
    purge-homeassistant-statistics-database
fi

if bashio::config.true "retention.enabled" && bashio::config.equals "retention.periodic" "daily"; then
    export-all-databases
fi
