# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-backup/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB-backup service
# ==============================================================================

s6-mkfifodir ${ADDON_EVENTS}

while :
do
    s6-ftrig-wait ${ADDON_EVENTS} "e" > /dev/null
    export-all-databases
    s6-ftrig-notify ${ADDON_EVENTS} "E"
done
