# mariadb/rootfs/etc/s6-overlay/s6-rc.d/cont/run
#!/command/with-contenv bashio
# ==============================================================================
# Prepare the MariaDB add-on for running
# ==============================================================================
TMPFS_MNT=/tmp
ADDON_EVENTS=${TMPFS_MNT}/events
MARIADB_DATA=${TMPFS_MNT}/databases
MARIADB_DUMP=/data/homeassistant-database-dump-TYPE.sql.gz
MARIADB_DUMP_FINISHED=/data/homeassistant-database-dump-finished
STATISTICS_DUMP=/share/purged-statistics/purged-homeassistant-statistics-database-older-than-TIMESTAMP.sql.gz
EXPORT_IMPORT_DATABASE_CONTENT=false

# Save variables
echo ${TMPFS_MNT} > /etc/TMPFS_MNT
echo ${ADDON_EVENTS} > /etc/ADDON_EVENTS
echo ${MARIADB_DATA} > /etc/MARIADB_DATA
echo ${MARIADB_DUMP} > /etc/MARIADB_DUMP
echo ${MARIADB_DUMP_FINISHED} > /etc/MARIADB_DUMP_FINISHED
echo ${STATISTICS_DUMP} > /etc/STATISTICS_DUMP
echo ${EXPORT_IMPORT_DATABASE_CONTENT} > /etc/EXPORT_IMPORT_DATABASE_CONTENT
