# mariadb/rootfs/etc/s6-overlay/s6-rc.d/cont/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Prepare the MariaDB add-on for running
# ==============================================================================
TMPFS_MNT=/tmp
MARIADB_DATA=${TMPFS_MNT}/databases
OLD_MARIADB_DUMP=/data/homeassistant-database-dump-TYPE.sql.gz
OLD_MARIADB_DUMP_FINISHED=/data/homeassistant-database-dump-finished
MARIADB_DUMP=/data/databases-dump.sql.gz
MARIADB_DUMP_FINISHED=/data/databases-dump-finished
STATISTICS_DUMP=/share/purged-statistics/purged-homeassistant-statistics-database-older-than-TIMESTAMP.sql.gz
EXPORT_IMPORT_ALL_DATABASES_CONTENT=false

# Save variables
printf "${TMPFS_MNT}" > /var/run/s6/container_environment/TMPFS_MNT
printf "${MARIADB_DATA}" > /var/run/s6/container_environment/MARIADB_DATA
printf "${OLD_MARIADB_DUMP}" > /var/run/s6/container_environment/OLD_MARIADB_DUMP
printf "${OLD_MARIADB_DUMP_FINISHED}" > /var/run/s6/container_environment/OLD_MARIADB_DUMP_FINISHED
printf "${MARIADB_DUMP}" > /var/run/s6/container_environment/MARIADB_DUMP
printf "${MARIADB_DUMP_FINISHED}" > /var/run/s6/container_environment/MARIADB_DUMP_FINISHED
printf "${STATISTICS_DUMP}" > /var/run/s6/container_environment/STATISTICS_DUMP
printf "${EXPORT_IMPORT_ALL_DATABASES_CONTENT}" > /var/run/s6/container_environment/EXPORT_IMPORT_ALL_DATABASES_CONTENT

# Copy scripts to share folder
mkdir -p /share/bin || true
cp -pf /usr/bin/update-mariadb-inmemory-while-ha-core-restarted /share/bin/ || true
cp -pf /usr/bin/restart-mariadb-inmemory-while-ha-core-restarted /share/bin/ || true
