# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-export/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Export all databases for backups
# ==============================================================================

# Export all databases content
export-all-databases
