# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start MariaDB service
# ==============================================================================

# Start mariadb
bashio::log.info "Starting MariaDB"
mkdir -p /run/mariadbd
exec mariadbd --datadir="${MARIADB_DATA}" --user=root < /dev/null
