# mariadb/rootfs/etc/periodic/hourly/hourly_export
#!/command/with-contenv bashio
# ==============================================================================
# Export homeassistant database periodically
# ==============================================================================

if bashio::config.true "retention.enabled" && bashio::config.equals "retention.periodic" "hourly"; then
    export-databases
fi
