# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core-pre/run
#!/command/with-contenv bashio
# ==============================================================================
# Pre-start initialization of MariaDB service
# ==============================================================================

# Change tmpfs size to configured value
TMPFS_SIZE=$(bashio::config "tmpfs.size")
bashio::log.info "Change tmpfs size to ${TMPFS_SIZE}"
mount -o remount,size=${TMPFS_SIZE} tmpfs ${TMPFS_MNT}

# Redirect log output
mkdir -p "${MARIADB_DATA}"
ln -s /proc/1/fd/1 "${MARIADB_DATA}/mariadb.err"

# Init mariadb
bashio::log.info "Create a new mariadb initial system"
mysql_install_db --user=root --datadir="${MARIADB_DATA}" --skip-name-resolve --skip-test-db > /dev/null
