# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb-core-pre/up
/etc/s6-overlay/scripts/mariadb-core-pre-up
