# mariadb/rootfs/etc/periodic/weekly/weekly_export
#!/command/with-contenv bashio
# ==============================================================================
# Export homeassistant database periodically
# ==============================================================================

if bashio::config.true "retention.enabled" && bashio::config.has_value "retention.periodic" "weekly"; then
    export-homeassistant-database
fi
